# debian/postinst
#!/bin/sh
set -e

addgroup --system scmcauthd
addgroup --system scmcshadow

adduser --system --home=/run/scmcauthd --no-create-home --ingroup scmcauthd scmcauthd
usermod scmcauthd -g scmcauthd -aG scmcshadow

/usr/lib/scmcauthd/create_run_directory
systemctl enable scmcauthd.service || true
systemctl start scmcauthd.service || true
