# debian/scmcauthd.prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ] || [ "$1" = "purge" ]
then
  systemctl stop scmcauthd.service || true
  systemctl disable scmcauthd.service || true
fi
